# example/x509-cert-expiry-pager.timer
Description=x509 cert expiry pager

[Timer]
OnCalendar=*-*-* 11:00:00 UTC
